# deployment/digital-twin.service
[Unit]
Description=Agent X Digital Twin
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/agent-x/digital-twin
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=append:%h/agent-x/logs/digital-twin.log
StandardError=append:%h/agent-x/logs/digital-twin.log

[Install]
WantedBy=multi-user.target
